# misc/deepin-service-manager.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service

[Service]
ExecStart=/usr/bin/deepin-service-manager
Restart=always
RestartSec=3
@SYSTEMD_SLICE@
# 由于服务区分root和user，安全加固在CMakelist.txt中区分设置
@SYSTEMD_HARDENING@

@SYSTEMD_INSTALL@
