# misc/deepin-service-plugin@.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service

[Service]
ExecStart=/usr/bin/deepin-service-manager -n %i
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@

ProtectSystem=strict
InaccessiblePaths=/etc/shadow
ReadWritePaths=-/usr/share/deepin-service-manager/
ReadWritePaths=-/usr/local/lib/deepin-service-manager/
ReadWritePaths=-/usr/local/deepin-service-manager/
