# misc/deepin-service-group@.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service deepin-service-manager.service
PartOf=deepin-service-manager.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart=/usr/bin/deepin-service-manager -g %i
Restart=on-failure
RestartSec=3
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@
