# misc/deepin-service-manager.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart=/usr/bin/deepin-service-manager
Restart=on-failure
RestartSec=3
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@
