# misc/deepin-service-plugin@.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service @EXTRA_AFTER@

[Service]
ExecStart=/usr/bin/deepin-service-manager -n %i
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@
